# SocketProgramming/echo.service
[Unit]
Description=Demo C service

[Service]
ExecStart=/home/ps/SocketProgramming/echoserver

[Install]
WantedBy=multi-user.target